# tools/whisper-stt.service
[Unit]
# Self-hosted whisper.cpp STT endpoint for the soledex Pokedex Tab5.
# Serves POST /inference (WAV in -> transcript out).
#
# Bound to localhost only and exposed to the internet via nginx (HTTPS + basic
# auth) -- see the "Hosting the STT server" section of the README. For a
# LAN-only box with no nginx in front, change --host to 0.0.0.0 and firewall
# port 8080 to the local subnet instead.
Description=soledex whisper.cpp speech-to-text server
After=network-online.target
Wants=network-online.target

[Service]
# Adjust User/paths to match the server. WHISPER_HOME is where you cloned
# https://github.com/ggml-org/whisper.cpp and built whisper-server.
User=whisper
Environment=WHISPER_HOME=/opt/whisper.cpp
# The whisper-server binary links libwhisper/libggml from build/bin.
Environment=LD_LIBRARY_PATH=/opt/whisper.cpp/build/bin
WorkingDirectory=/opt/whisper.cpp
ExecStart=/opt/whisper.cpp/build/bin/whisper-server \
  -m /opt/whisper.cpp/models/ggml-base.en.bin \
  -l en --host 127.0.0.1 --port 8080
Restart=on-failure
RestartSec=2
# Hardening (relax if the model lives outside the home above).
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadOnlyPaths=/opt/whisper.cpp
PrivateTmp=true

[Install]
WantedBy=multi-user.target
